# opencode/.config/systemd/user/opencode.service
[Unit]
Description=OpenCode Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=OPENCODE_EXPERIMENTAL=1
ExecStart=%h/.bin/opencode-serve --hostname 0.0.0.0 --port 4096
LoadCredential=OPENCODE_SERVER_PASSWORD:%h/.credentials/opencode_password
LoadCredential=SUPERMEMORY_API_KEY:%h/.credentials/supermemory_api_key
Restart=on-failure
RestartSec=10

# Security hardening (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
